# src/templates/template.service
[Unit]
Description=$description

[Service]
WorkingDirectory=$working_dir
ExecStart=$command
