# src/templates/template.timer
[Unit]
Description=$description
$unit_options

[Timer]
$timer_action

[Install]
WantedBy=timers.target
