# src/kshift/templates/template.service
[Unit]
Description=$description

[Service]
Environment="SOURCE=systemd"
ExecStart=$command
